# deploy/systemd/cleo-rover-senses.service
# Example unit for the offline voice daemon (wake word -> STT -> /pip/command).
# Requires the `voice` extra installed plus an offline STT backend (whisper.cpp or
# vosk); see rover/voice_daemon.py. Adjust User=/paths to match your install; the
# body service must be running. (For the old low-rate event simulator instead, use
# ExecStart=.../cleo-rover-senses --simulate --interval 10.)
[Unit]
Description=Cleo Rover voice daemon
After=cleo-rover-body.service
Requires=cleo-rover-body.service

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/cleo-rover
ExecStart=/home/pi/cleo-rover/.venv/bin/cleo-rover-voice --base http://127.0.0.1:8099
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
